# Misa-Cash/src/web/backend/scripts/crontab
# Backup do banco de dados (diário à meia-noite)
0 0 * * * /usr/local/bin/python3 /app/scripts/backup.py >> /var/log/misacash/backup.log 2>&1

# Limpeza de logs antigos (semanal aos domingos à 1h)
0 1 * * 0 find /var/log/misacash -name "*.log" -mtime +30 -exec rm {} \;

# Verificação de integridade do banco (diário às 2h)
0 2 * * * /usr/local/bin/python3 /app/scripts/check_db_integrity.py >> /var/log/misacash/db_check.log 2>&1

# Geração de relatórios diários (diário às 3h)
0 3 * * * /usr/local/bin/python3 /app/scripts/generate_daily_reports.py >> /var/log/misacash/reports.log 2>&1

# Limpeza de cache (diário às 4h)
0 4 * * * /usr/local/bin/python3 /app/scripts/clear_cache.py >> /var/log/misacash/cache.log 2>&1

# Verificação de segurança (semanal às segundas às 5h)
0 5 * * 1 /usr/local/bin/python3 /app/scripts/security_check.py >> /var/log/misacash/security.log 2>&1

# Backup de arquivos de mídia (diário às 6h)
0 6 * * * /usr/local/bin/python3 /app/scripts/backup_media.py >> /var/log/misacash/media_backup.log 2>&1

# Verificação de performance (diário às 7h)
0 7 * * * /usr/local/bin/python3 /app/scripts/check_performance.py >> /var/log/misacash/performance.log 2>&1

# Manutenção do banco de dados (semanal aos sábados às 8h)
0 8 * * 6 /usr/local/bin/python3 /app/scripts/db_maintenance.py >> /var/log/misacash/db_maintenance.log 2>&1

# Verificação de backups (diário às 9h)
0 9 * * * /usr/local/bin/python3 /app/scripts/verify_backups.py >> /var/log/misacash/backup_verify.log 2>&1 